# chef/cookbooks/couchdb/templates/default/couchdb.build.conf.erb
description "<%= @couch_service_name %>"
start on filesystem
stop on runlevel S
respawn

exec <%= @cwd %>/build/bin/couchdb


